# scripts/rubygem_secret_monitor.timer
[Unit]
<<<<<<< HEAD
Description=Run rubygem_secret_monitor every day at 7am pacific
Requires=rubygem_secret_monitor.service

[Timer]
OnCalendar=*-*-* 15:00:00
RandomizedDelaySec=600

[Install]
WantedBy=timers.target